# autosphinx.conf
[supervisord]
nodaemon=true

[program:sphinx]
command=/usr/bin/make html
directory=/usr/share/nginx/html
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
