# autosphinx
#!/bin/bash

### BEGIN INIT INFO
# Provides:          autosphinx
# Required-Start:    $all
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Initiates sphinx-autobuild server
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=/usr/sbin/nginx
NAME=autosphinx
DESC=autosphinx
PROJ=$(ls /sphinx/projects | tee /var/local/PROJ.txt)
PROJS=$(ls /sphinx/projects | wc -l)

STOP_SCHEDULE="${STOP_SCHEDULE:-QUIT/5/TERM/5/KILL/5}"

sphinx_auto () {
    if [[ "$PROJS" -gt "1" ]]
    then
        sphinx-autobuild --port=0 -a sphinx/default/source /sphinx/default/build/html
        while IFS= read -r LINE; do sphinx-autobuild --port=0 -a /sphinx/projects/$LINE/source /sphinx/projects/$LINE/build/html; done < /var/local/PROJ.txt
    elif [[ "$PROJS" -eq "1" ]]
    then
        sphinx-autobuild -a /sphinx/projects/$PROJ/source /sphinx/projects/$PROJ/build/html
    else
        sleep 0
    fi
}

# Try to extract nginx pidfile
PID=$(pgrep sphinx-autobuild*)

case "$1" in
    start) sphinx_auto &> /dev/null &
       ;;
    stop)
       kill $PID
       ;;
    restart)
       $0 stop
       $0 start
       ;;
    status)
       status_of_proc -p && exit 0 || exit $?
       ;;
    *)
       echo "Usage: $0 {start|stop|status|restart}"
esac
